# src/config/web-counter.conf
description	"Foo"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
chdir /root/web-counter
exec /root/web-counter/web-counter -p 80
